# server/ip_register/systemd/ip_register.timer
[Unit]
Description=timer ip_register

[Timer]
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target