# server/ip_register/systemd/ip_register.timer
[Unit]
Description=timer ip_register

[Timer]
OnUnitActiveSec=5m
Unit=ip_register.service

[Install]
WantedBy=timers.target